# profile.d/50-kabre.sh
#!/bin/sh

case "$HOSTNAME" in
  login-*.cnca) TERM=xterm ;;
esac
