# profile.d/50-nvm.sh
# shellcheck shell=sh

export NVM_DIR="${XDG_CONFIG_HOME:-$HOME/.config}"/nvm
