# profile.d/90-custom-vars.sh
# shellcheck shell=sh

export \
  EDITOR="nvim" \
  VISUAL="em" \
  BROWSER="firefox"
