# shell/profile.d/10-custom-paths.sh
# shellcheck shell=sh

export PATH="$HOME/.local/bin${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts${PATH:+:${PATH}}"
export PATH="$HOME/.config/wrappers${PATH:+:${PATH}}"
