# shell/profile.d/50-aspell.sh
# shellcheck shell=sh

export ASPELL_CONF="home-dir ${XDG_CONFIG_HOME:-$HOME/.config}/aspell"
