# shell/profile.d/50-cargo.sh
# shellcheck shell=sh

export PATH="${XDG_DATA_HOME:-$HOME/.local/share}/cargo/bin${PATH:+:${PATH}}"
