# shell/profile.d/50-doom.sh
# shellcheck shell=sh

export PATH="${XDG_CONFIG_HOME:-$HOME/.config}/emacs/bin${PATH:+:${PATH}}"
