# shell/profile.d/50-npm.sh
# shellcheck shell=sh

export PATH="${XDG_DATA_HOME:-$HOME/.local/share}/npm/bin${PATH:+:${PATH}}"
