# shell/profile.d/cargo.sh
#!/bin/sh

export PATH="$HOME/.local/share/cargo/bin:$PATH"
