# shell/profile.d/cargo.sh
# shellcheck shell=sh

export PATH="$HOME/.local/share/cargo/bin${PATH:+:${PATH}}"
