# shell/profile.d/config.sh
#export QT_QPA_PLATFORMTHEME="qt5ct"
#export GTK_THEME="Materia:dark"
#export XCURSOR_THEME="Adwaita"
export CM_SELECTIONS="clipboard"
export CM_MAX_CLIPS=10
