# shell/profile.d/config.sh
export \
    HISTCONTROL=ignoreboth \
    HISTIZE= \
    HISTFILESIZE=

export \
    QT_QPA_PLATFORMTHEME="qt5ct" \
    CM_SELECTIONS="clipboard" \
    GTK_THEME="Materia:dark" \
    XCURSOR_THEME="Adwaita" \
    CM_MAX_CLIPS=10
