# shell/profile.d/defaults.sh
export \
    TERM="xterm-256color" \
    BROWSER="firefox-hardened" \
    EDITOR="nvim" \
    VISUAL="nvim"
