# shell/profile.d/defaults.sh
export \
    TERM="xterm-256color" \
    BROWSER="firefox-hardened" \
    EDITOR="cvim" \
    VISUAL="cvim"
