# shell/profile.d/defaults.sh
#!/bin/sh

export \
  TERM="xterm-256color" \
  BROWSER="firefox" \
  EDITOR="lvim" \
  VISUAL="lvim"
