# shell/profile.d/doom.sh
# shellcheck shell=sh

export PATH="$HOME/.config/emacs/bin${PATH:+:${PATH}}"
