# shell/profile.d/go.sh
#!/bin/sh

export PATH="$HOME/.local/share/go/bin${PATH:+:${PATH}}"
