# shell/profile.d/nav.sh
#!/bin/sh

nav() {
    fzf \
    --cycle \
    --height=60% \
    --preview 'p {}' \
    --preview-window=right,70% \
    --prompt='⯈ ' \

}
