# shell/profile.d/nix.sh
# shellcheck shell=sh

export PATH="$HOME/.local/state/nix/profile/bin${PATH:+:${PATH}}"
