# shell/profile.d/path.sh
#!/bin/sh
# shellcheck disable=SC2123

export PATH="$HOME/.local/bin${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/dwm${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/menu${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/setup${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/status${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/sway${PATH:+:${PATH}}"
export PATH="$HOME/.config/scripts/x11${PATH:+:${PATH}}"
export PATH="$HOME/.config/wrappers${PATH:+:${PATH}}"
