# shell/profile.d/pyenv.sh
# shellcheck shell=sh

export PATH="$HOME/.local/share/pyenv/shims${PATH:+:${PATH}}"
