# shell/profile.d/shortcuts.sh
export \
    GIT="git@git.tavo.one:tavo" \
    GITLAB="git@gitlab.com:tavo-wasd" \
    GITHUB="git@github.com:tavo-wasd-gh" \
    DRIVE="ssh://drive:/home/drive/drive" \
    BOOKMARKS="$HOME/Documents/bookmarks" \
    BIB="$HOME/Documents/bibliography" \
    DATE=$(date -I) \
    WEEK=$(date '+%U')
