# systemd/system/getty@tty1.service.d/override.conf
# Place this file in:
# /etc/systemd/system/getty@tty1.service.d/override.conf
# to autologin

[Service]
ExecStart=
ExecStart=-/sbin/agetty -a tavo %I
